# v2ray-linux-amd64/DEBIAN/postrm
#!/bin/bash

systemctl stop v2ray
systemctl disable v2ray
rm -rf /var/log/v2ray
