# v2ray-linux-amd64/DEBIAN/prerm
#!/bin/bash

unlink /etc/v2ray/config.json
systemctl stop v2ray
systemctl disable v2ray
